# files/usr/var/service/httpd/run
#!/data/data/com.termux/files/usr/bin/sh
exec httpd -DNO_DETACH 2>&1
