# files/usr/var/service/postgres/run
#!/data/data/com.termux/files/usr/bin/sh
mkdir -p ~/.postgres
exec postgres -D ~/.postgres/ 2>&1
